# Xtra.ini
[program:uwsgi_django]
command=/root/.virtualenvs/Xtra/bin/uwsgi --ini /root/Xtra/XtrawithDjango/Xtra/wsgi.ini
directory=/root/Xtra/XtrawithDjango/Xtra
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true

[supervisord]
nodaemon=false

[supervisorctl]
